# scripts/tcpcheck_iptable/tcpcheck.service
[Unit]
After=tcpcheck.service

[Service]
ExecStart=/etc/tcpcheck.sh

[Install]
WantedBy=default.target
